# skel/etc/supervisor/conf.d/300-zookeeper.conf
[program:zookeeper]
command=%(ENV_SERVICE_ZOOKEEPER_CMD)s
user=root
autostart=true
autorestart=unexpected
startsecs=5
startretries=10
priority=300
stopsignal=KILL
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
